# .skills/openclaw-skills/skills/x3r081/blocket-watch/systemd/blocket-watch.timer
[Unit]
Description=Run Blocket watch every 60 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
